# com.bedichek.watch_gobble_solartracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" 
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bedichek.watch_gobble_tracker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/robertbedichek/bin/watch_serial_and_run_gobble_tracker.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/watch_serial.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/watch_serial.err</string>
</dict>
</plist>

